# .docker/services/webservice/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command = /usr/sbin/php-fpm8.1 --nodaemonize -c /etc/php/8.1/fpm/php-fpm.conf
socket = unix:///run/php/php-fpm.sock
autorestart = true
startretries = 5
priority = -2
user = root
